# custom-files/etc/sv/nvidia-persistenced/run
#!/bin/sh
nvidia-smi -pm 1
sv down nvidia-persistenced
